# roles/zfs_remote_backup/templates/zfs-remote-backup.service.j2
[Unit]
Description=ZFS Remote Backup Job

[Service]
Type=oneshot
ExecStart=/usr/local/bin/zfs-remote-backup.sh
